# systemd/bookmarks.service
[Unit]
Description=Slack Bookmarks Flask Service
After=network.target

[Service]
Type=simple
User=goce
WorkingDirectory=/mnt/ssd/apps/bookmarks
# Load webhook URL and token from .env (copy from repo apps/bookmarks/.env.example)
EnvironmentFile=-/mnt/ssd/apps/bookmarks/.env
Environment="PATH=/mnt/ssd/apps/bookmarks/venv/bin:/usr/local/bin:/usr/bin:/bin:/home/goce/.local/bin"
ExecStart=/mnt/ssd/apps/bookmarks/venv/bin/python3 /mnt/ssd/apps/bookmarks/secure_slack_bookmarks.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

