# systemd/knowledge-mcp-weekly-refresh.service
[Unit]
Description=Weekly MCP knowledge base refresh (Jira export + ingest)
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=goce
Group=goce
WorkingDirectory=/home/goce/Desktop/Cursor projects/mcp_server
ExecStart=/bin/bash /home/goce/Desktop/Cursor projects/mcp_server/scripts/weekly-knowledge-refresh.sh
# Jira export + embeddings can take 30–60+ minutes
TimeoutStartSec=7200

[Install]
WantedBy=multi-user.target
