# systemd/portfolio-update.timer
[Unit]
Description=Run Portfolio Auto-Update Every 5 Minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
