# systemd/portfolio-update.timer
[Unit]
Description=Portfolio Auto-Update Timer (DISABLED - Use 'make portfolio-update' instead)
# This timer is disabled to reduce CPU usage. Use 'make portfolio-update' to manually update.

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
